# com.user.obsidian-summarizer.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!-- Example launchd configuration file for Obsidian Summarizer -->
<!-- Copy this file to ~/Library/LaunchAgents/com.user.obsidian-summarizer.plist -->
<!-- Replace [USERNAME] and [PATH] with your actual values -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.obsidian-summarizer</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/[USERNAME]/path/to/obsidian-summarizer/run_summarizer.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/[USERNAME]/path/to/obsidian-summarizer/cron.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/[USERNAME]/path/to/obsidian-summarizer/cron.log</string>
    <key>RunAtLoad</key>
    <true/>
    <key>StartOnMount</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
</dict>
</plist> 